# pkg/checkpoint/DEBIAN/prerm
#!/bin/sh
set -e

echo "Stopping Checkpoint service..."
systemctl stop checkpoint.service || true
systemctl disable checkpoint.service || true

echo "Checkpoint service stopped and disabled."
