# pkg/checkpoint/DEBIAN/prerm
#!/bin/sh
set -e
systemctl stop checkpoint.service || true
systemctl disable checkpoint.service || true
